# docker/file_system/etc/supervisor/conf.d/services.conf
[program:postgresql]
command=/bin/bash -lc 'exec "$$(command -v postgres)" -D /data/db/postgres'
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/log/postgres/postgresql.log
stderr_logfile=/data/log/postgres/postgresql.log

[program:valkey]
command=/usr/bin/valkey-server /config/valkey/valkey.conf
autostart=true
autorestart=true
priority=20
stdout_logfile=/data/log/valkey.log
stderr_logfile=/data/log/valkey.log

[program:tor]
command=/usr/bin/tor -f /config/caspbx/tor/torrc
autostart=%(ENV_TOR_ENABLED)s
autorestart=true
priority=30
stdout_logfile=/data/log/tor.log
stderr_logfile=/data/log/tor.log

[program:app]
command=/usr/local/bin/caspbx --address 0.0.0.0 --port 80
autostart=true
autorestart=true
priority=100
stdout_logfile=/data/log/app.log
stderr_logfile=/data/log/app.log
